# email_sender/email-sender.service.example
[Unit]
Description=LibertyCall Email Sender Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/libertycall
Environment="PATH=/opt/libertycall/venv/bin"
ExecStart=/opt/libertycall/venv/bin/python -m email_sender.scheduler_service
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

