# deploy/systemd/gateway.service
[Unit]
Description=LibertyCall Gateway Service
After=network.target

[Service]
Type=simple
# User=liberty  # libertyユーザーが存在しないためコメントアウト
WorkingDirectory=/opt/libertycall
ExecStart=/opt/libertycall/venv/bin/python3 /opt/libertycall/gateway/realtime_gateway.py --rtp_port 7002
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
# WatchdogSec=30  # watchdog timeoutを無効化（Gatewayが正常に動作するまで）
Environment="PYTHONUNBUFFERED=1"
StandardOutput=append:/opt/libertycall/logs/gateway_stdout.log
StandardError=append:/opt/libertycall/logs/gateway_stderr.log

[Install]
WantedBy=multi-user.target
